# docker/files/pueued.service
[Unit]
Description=Pueued

[Service]
Restart=yes
ExecStart=/usr/bin/pueued -vv

[Install]
WantedBy=default.target
